# deploy/lucy.service
[Unit]
Description=Lucy Slack Bot (AI Coworker)
Documentation=https://github.com/ojashyadav101/lucy
After=network.target docker.service nginx.service
Wants=docker.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lucy/app
EnvironmentFile=/opt/lucy/app/.env

ExecStart=/opt/lucy/app/.venv/bin/uvicorn lucy.app:api --host 127.0.0.1 --port 3000 --workers 1

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lucy

# Allow in-flight agent runs to drain on restart
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
